# docker/php/supervisor/supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0766                       ; sockef file mode (default 0700)

[supervisord]
nodaemon = true
directory = /tmp
childlogdir = /tmp
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
logfile_maxbytes = 50MB
logfile_backups = 10

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf

[program:php-fpm]
command = /usr/bin/php-fpm -F
priority=100
autostart = true
autorestart = unexpected
stdout_logfile = /tmp/php-stdout.log
stdout_logfile_maxbytes = 0
stderr_logfile = /tmp/php-stderr.log
stderr_logfile_maxbytes = 0
exitcodes = 0

; Laravel queue workers
[program:laravel-queue]
process_name = %(program_name)s_%(process_num)02d
command = php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart = true
autorestart = unexpected
stopasgroup = true
killasgroup = true
user = www-data
numprocs = 1
redirect_stderr = true
stdout_logfile = /tmp/queue.log
stopwaitsecs = 3600
